# install/daemons/enco_queue.conf
start on startup
stop on shutdown
exec /enco/local/core/start_queue_daemon.sh
respawn
