# install/daemons/rali_calendar.conf
start on startup
stop on shutdown
exec /rali/rali-app/local/core/start_calendar_daemon.sh
respawn
